# erouter/config/services.d/erouter-core/run
#!/usr/bin/with-contenv sh

LOG_FILE=/home/erouter/erouter/tmp/e-router_core.log
cd /home/erouter/erouter \
&& /home/erouter/erouter/e-router start -t -- core >>$LOG_FILE 2>&1
