# erouter/config/services.d/erouter-veiled.oracle/run
#!/usr/bin/with-contenv sh

LOG_FILE=/home/erouter/erouter/tmp/veiled_oracle.log
cd /home/erouter/erouter \
&& DEBUG=1 /home/erouter/erouter/veiled_controller run oracle >>$LOG_FILE 2>&1

